# multiple-network-interface-fixes/rootfs/etc/s6-overlay/s6-rc.d/mnifix/run
#!/usr/bin/with-contenv bashio
set -e

# Start the service
exec /usr/bin/mnifix-run
